# nft-listener.service
[Unit]
Description=NFT Firewall Listener — Keybase ChatOps bot
After=network-online.target nft-watchdog.service
Wants=network-online.target
Requires=nft-watchdog.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/fw-admin/nft-firewall/src/main.py listener daemon
Restart=always
RestartSec=15
User=root
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nft-listener

[Install]
WantedBy=multi-user.target
