# systemd/nft-firewall-doctor.service
[Unit]
Description=NFT Firewall doctor check

[Service]
Type=oneshot
ExecStart=/usr/local/bin/fw doctor
